# files/conf/rc.local
#!/bin/sh -e

su - display -c startx
exit 0
